# crontab
* * * * * /var/www/executori.org/date/backup.sh
