# crontab
0 */3 * * * /var/www/executori.org/date/backup.sh
